# tval.service
[Unit]
Description=TVAL

[Service]
Type=simple
ExecStart=sudo python3 /home/pi/TVAL/main.py

[Install]
WantedBy=multi-user.target